# scripts/acdc/launchd/com.bason.acdc-watch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bason.acdc-watch</string>

  <!-- launchd provides only a minimal PATH (/usr/bin:/bin:...) and does NOT source
       shell rc files, so node/pnpm/gh/claude are not found by default. install.sh
       detects their real locations and bakes them into PATH (__PATH__). We use a
       NON-login `zsh -c` because the login shell (-l) tripped a getcwd error under
       launchd while sourcing rc files. __HOME__/__REPO__/__PATH__ are substituted
       by install.sh. The watcher self-loads ~/.acdc/board.env, claude-token.env,
       and firebase.env, so no other env wiring is needed here. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>cd "__REPO__" &amp;&amp; exec pnpm -C scripts/acdc exec tsx bin/acdc-watch.ts</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/.acdc/watch.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/.acdc/watch.log</string>
</dict>
</plist>
